# molecule/default/files/borgmatic.timer
# Ansible managed: Do NOT edit this file manually!

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=6h

[Unit]
Description=Run borgmatic backup
